# openpower-pnor-update@.service
[Unit]
Description=Update PNOR %I

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/sbin/pflash -E -f -p %I
SyslogIdentifier=pflash

